# systemd/api-transformer.socket
[Unit]
Description=api-transformer socket

[Socket]
ListenStream=/run/gunicorn/api-transformer.socket

[Install]
WantedBy=sockets.target
